# deploy/bilive-smb-recover.timer
[Unit]
Description=Periodically recover bilive Windows SMB mount

[Timer]
OnBootSec=10s
OnUnitActiveSec=15s
AccuracySec=1s
Unit=bilive-smb-recover.service

[Install]
WantedBy=timers.target

