# deploy/bilive-dashboard.service
[Unit]
Description=Bilive slice dashboard
After=network-online.target bilive.service
Wants=network-online.target
Requires=bilive.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/mnt/win/bilive
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
Environment=BILIVE_VIDEOS_DIR=/mnt/win/bilive/Videos
EnvironmentFile=-/mnt/win/bilive/.secrets/env
ExecStart=/mnt/win/bilive/start_dashboard_pi.sh
Restart=on-failure
RestartSec=10
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
