# deploy/bilive.service
[Unit]
Description=Bilive B站直播录制
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
EnvironmentFile=/mnt/win/bilive/.secrets/env
Environment=BILIVE_VIDEOS_DIR=/mnt/win/bilive/Videos
Environment=BILIVE_LOG_DIR=/mnt/win/bilive/logs
Environment=BILIVE_COOKIE_FILE=/mnt/win/bilive/.secrets/bilibili.cookie
ExecStart=/usr/local/bin/bilive-start.sh
Restart=on-failure
RestartSec=10
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target