# deploy/bilive-dashboard.service
[Unit]
Description=Bilive slice dashboard
After=network-online.target mnt-win.automount
Wants=network-online.target mnt-win.automount bilive-smb-recover.timer

[Service]
Type=simple
User=ubuntu
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
Environment=BILIVE_VIDEOS_DIR=/mnt/win/bilive/Videos
Environment=BILIVE_WAIT_SECONDS=15
Environment=BILIVE_HEALTHCHECK_SECONDS=15
ExecStart=/usr/local/bin/bilive-dashboard-start.sh
Restart=always
RestartSec=10
TimeoutStopSec=30
MemoryAccounting=yes
MemoryMax=512M
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
